# templates/iperf.conf
description "IPERF Daemon"

start on runlevel[2345]
stop on runlevel [!2345]

exec iperf -s

